# supervisord.conf
[unix_http_server]
file=/var/run/supervisor.sock

[supervisord]
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log
nodaemon=true

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[program:sshd]
command=/usr/sbin/sshd -D
autostart=true
autorestart=true
user=root
stdout_logfile=/var/log/sshd.out.log
stderr_logfile=/var/log/sshd.err.log

[program:xray]
command=/usr/local/bin/xray -config /etc/xray/config.json
autostart=true
autorestart=true
user=root
stdout_logfile=/var/log/xray.out.log
stderr_logfile=/var/log/xray.err.log

[program:sing-box]
command=/usr/local/bin/sing-box run -c /etc/sing-box/config.json
autostart=true
autorestart=true
user=root
stdout_logfile=/var/log/sing-box.out.log
stderr_logfile=/var/log/sing-box.err.log

[program:openlist]
command=/usr/local/bin/openlist server
autostart=true
autorestart=true
user=root
stdout_logfile=/var/log/openlist.out.log
stderr_logfile=/var/log/openlist.err.log

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
user=root
stdout_logfile=/var/log/nginx_supervisor.log
stderr_logfile=/var/log/nginx_supervisor_error.log
